# jupyter.ini
[program:jupyter]
command=pixi run jupyterhub
directory=/workspace
autostart=true
autorestart=true
stderr_logfile=/var/log/jupyter.err.log
stdout_logfile=/var/log/jupyter.out.log
startsecs=5
