# config/setpassword.ini
[program:setpassword]
command=/usr/local/bin/setpassword.sh
directory=/workspace
autostart=true
autorestart=false
exitcodes=0
startsecs=0
priority=1

